# systemd/gluetun-ap-routing.service
[Unit]
Description=Routing and iptables for AP clients via gluetun
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/gluetun-ap/apply-routing.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
