# systemd/gluetun-ap-routing.service
[Unit]
Description=Routing and iptables for AP clients via gluetun
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
ExecStart=/usr/local/lib/gluetun-ap/apply-routing.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
