# systemd/gluetun-ap-firewall.service
[Unit]
Description=Firewall rules to route AP clients through gluetun and block WAN leaks
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/gluetun-ap/apply-firewall.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
